# infrastructure/systemd/superparty-worker-ads.service
[Unit]
Description=Superparty Ads Worker
After=network.target

[Service]
Type=simple
User=superparty-agent
Group=superparty-agent
WorkingDirectory=/var/www/superparty
EnvironmentFile=/var/www/superparty/.env
ExecStart=/usr/bin/python3 -m agent.services.worker ads_plan,ads_draft
Restart=always
RestartSec=5
StandardOutput=append:/var/www/superparty/logs/worker_ads.log
StandardError=append:/var/www/superparty/logs/worker_ads.log

[Install]
WantedBy=multi-user.target
